# consul_config/consul_config.cron
# 备份config 文件
10 16 * * * /bin/sh /app/check.sh
